# supervisord.conf
[supervisord]
nodaemon=true
[supervisorctl]

[program:sshd]
command=/usr/sbin/sshd -D
[program:nginx]
command=/usr/sbin/nginx
[include]
files = /home/damir/data/conf/*.ini