# data/supervisord.conf
[supervisord]
nodaemon=true

[include]
files = /etc/supervisor/conf.d/*.conf

[program:sshd]
command=/usr/sbin/sshd -D
autorstart=true
autorestart=false
[program:nginx]
command=/usr/sbin/nginx
autorstart=true
autorestart=false
[program:teamview]
command=/home/damir/data/run/Project/Web/teamview/teamview -importPath teamview -srcPath "/home/damir/data/run/Project/Web/teamview/src" -runMode prod
autorstart=false
directory=/home/damir/data/run/Project/Web/teamview
autorestart=false
startsecs=10
startretries=36
stdout_logfile=/home/damir/data/run/Project/Web/teamview/app.log
stderr_logfile=/home/damir/data/run/Project/Web/teamview/app.err